# TOS 7-template-deb-single/DEBIAN/postinst
#!/bin/bash
# Post-installation script — runs after the package is installed
set -e

APPID="myapp"
APP_DIR="/usr/local/${APPID}"

# Create dedicated user
if ! id -u "${APPID}" > /dev/null 2>&1; then
    useradd --system --no-create-home --shell /usr/sbin/nologin "${APPID}" 2>/dev/null || true
fi

# Set permissions
chown -R "${APPID}:${APPID}" "${APP_DIR}" 2>/dev/null || true
chown -R "${APPID}:${APPID}" "/var/lib/${APPID}" 2>/dev/null || true

# Decompress webui.bz2 if present
if [ -f "${APP_DIR}/webui.bz2" ]; then
    (cd "${APP_DIR}" && tar -xjf webui.bz2 2>/dev/null) || true
fi

# Enable and start service
systemctl daemon-reload
systemctl enable "${APPID}.service" 2>/dev/null || true
systemctl start "${APPID}.service" 2>/dev/null || true

echo "${APPID} installed successfully."
exit 0
