# TOS 7-template-deb-dual/data-pkg/DEBIAN/prerm
#!/bin/bash
set -e

APPID="myapp"

systemctl stop "${APPID}.service" 2>/dev/null || true
systemctl disable "${APPID}.service" 2>/dev/null || true

exit 0
